# mainrig-MS-7E49/.config/systemd/user/psd-resync.timer.d/frequency.conf
[Unit]
Description=Timer for Profile-sync-daemon - 10min
BindsTo=

[Timer]
OnCalendar=
OnUnitActiveSec=
OnBootSec=5min
OnUnitActiveSec=10min
