# mainrig-MS-7E49/.config/systemd/user/psd-resync.service.d/override.conf
[Unit]
Description=Timer for Profile-sync-daemon - 10min

[Timer]
OnUnitActiveSec=
OnActiveSec=2min
OnUnitActiveSec=10min
