# mainrig-MS-7E49/.config/systemd/user/psd-resync.service.d/override.conf
[Service]
StandardOutput=null
StandardError=null
