# mainrig-MS-7E49/.config/systemd/user/llama-swap.service
# ~/.config/systemd/user/llama-swap.service
# Enable with:  systemctl --user enable --now llama-swap
# (blu already has `loginctl enable-linger` per your sync.sh, so this
#  survives with no active session, same as your commit job.)

[Unit]
Description=llama-swap — OpenAI endpoint that hot-swaps local models
After=graphical-session.target

[Service]
ExecStart=/usr/bin/llama-swap --listen 127.0.0.1:11500 --config %h/.config/llama-swap/config.yaml
Restart=on-failure
RestartSec=3
# A 60+ GB RAM-resident model shouldn't be first against the wall under
# memory pressure; nudge the OOM killer toward other things first.
OOMScoreAdjust=-200

[Install]
WantedBy=default.target
