# mainrig-MS-7E49/.config/systemd/user/refresh-marks.path
[Unit]
Description=Monitor Index.yml for changes

[Path]
PathModified=/etc/nginx/Index.yml

[Install]
WantedBy=default.target
