# mainrig-MS-7E49/.config/systemd/user/refresh-marks.service
[Unit]
Description=Rebuild static-marks index

[Service]
Type=oneshot
# This ensures systemd can find 'node' and other tools
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
WorkingDirectory=/etc/nginx
# Use the full path for both the shell and the redirect
ExecStart=/bin/sh -c '/usr/bin/static-marks build /etc/nginx/Index.yml > /etc/nginx/index.html'
