# mainrig-MS-7E49/.config/systemd/user/llama-swap.service.d/10-path.conf
[Service]
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/home/blu/.local/bin:/home/blu/.local/share/soar/bin
